# debian/dropbear-lite.postrm
#!/bin/sh
set -e

test "$1" = 'purge' || exit 0
if test -e /etc/dropbear; then
  rm -rf /etc/dropbear
fi
update-rc.d dropbear-lite remove >/dev/null
